# src/debian/config/prerm
#!/bin/sh

if [ -x "/etc/init.d/glassfish3" ]; then
  invoke-rc.d glassfish3 stop || true
fi
