# src/debian/bin/glassfish3
#!/bin/sh

#
# /etc/init.d/glassfish3 -- startup script for the JavaEE Glassfish3 server
#
### BEGIN INIT INFO
# Provides:          glassfish3
# Required-Start:    $local_fs $remote_fs $network
# Required-Stop:     $local_fs $remote_fs $network
# Should-Start:      $named
# Should-Stop:       $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start Glassfish.
# Description:       Start the JavaEE Glassfish3 server
### END INIT INFO

set -e

GF3_HOME=/usr/share/glassfish3/glassfish
DOMAINS_HOME=/var/lib/glassfish3
USER=glassfish

LOGS_DIR=/var/log/glassfish3

if [ `id -u` -ne 0 ]; then
        echo "You need root privileges to run this script"
        exit 1
fi

# Make sure tomcat is started with system locale
if [ -r /etc/default/locale ]; then
        . /etc/default/locale
        export LANG
fi

. /lib/lsb/init-functions

if [ -r /etc/default/rcS ]; then
        . /etc/default/rcS
fi

# load the glassfish3 configuration file
. /etc/default/glassfish3

##################################################
# Get the action & configs
##################################################
ACTION=$1

##################################################
# Setup JAVA if unset
##################################################
if [ -z "$JAVA" ]
then
  JAVA=$(which java)
fi

if [ -z "$JAVA" ]
then
  echo "Cannot find a Java JDK. Please set either set JAVA or put java (>=1.6) in your PATH." 2>&2
  exit 1
fi

##################################################
# Do the action
##################################################
case "$ACTION" in

  start)
    echo "Starting Glassfish please wait ....."
    sudo -u $USER $GF3_HOME/bin/asadmin --user $ADMIN_USER --passwordfile $PASSWORD_FILE start-domain $DOMAIN_NAME

    exit 1
    ;;

  stop)
    echo "Stopping Glassfish please wait ....."
    sudo -u $USER $GF3_HOME/bin/asadmin --user $ADMIN_USER --passwordfile $PASSWORD_FILE stop-domain $DOMAIN_NAME
    exit 1
    ;;

  restart)
    echo "Restarting Glassfish please wait ....."
    sudo -u $USER $GF3_HOME/bin/asadmin --user $ADMIN_USER --passwordfile $PASSWORD_FILE restart-domain $DOMAIN_NAME

    exit 1
    ;;

  create-domain)
    sudo -u $USER $GF3_HOME/bin/asadmin --user $ADMIN_USER --passwordfile $PASSWORD_FILE create-domain --adminport $ADMIN_PORT --savemasterpassword=true --savelogin=false --checkports=true --nopassword=false $DOMAIN_NAME

    sudo rm -rf $DOMAINS_HOME/$DOMAIN_NAME/logs
    sudo -u $USER mkdir -p $LOGS_DIR/$DOMAIN_NAME
    sudo -u $USER ln -s $LOGS_DIR/$DOMAIN_NAME $DOMAINS_HOME/$DOMAIN_NAME/logs

    sudo -u $USER $GF3_HOME/bin/asadmin --user $ADMIN_USER --passwordfile $PASSWORD_FILE start-domain $DOMAIN_NAME
    sudo -u $USER $GF3_HOME/bin/asadmin --user $ADMIN_USER --passwordfile $PASSWORD_FILE --host localhost --port $ADMIN_PORT enable-secure-admin
    sudo -u $USER $GF3_HOME/bin/asadmin --user $ADMIN_USER --passwordfile $PASSWORD_FILE stop-domain $DOMAIN_NAME

    exit 1
    ;;

  *)
    echo "usage: $0 (create-domain|start|stop|restart)"

    exit 1
    ;;
esac

exit 0
